# backend/supervisord.conf
[supervisord]
nodaemon=true
 
[group:server]
programs=tfrs,security-scan-script

[program:tfrs]
startsecs=30
startretries=20
redirect_stderr=true
directory=/app
command=bash -c "python3 manage.py runserver 0.0.0.0:8000"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:security-scan-script]
startsecs=10
startretries=20
redirect_stderr=true
directory=/app
command=bash -c "python3 manage.py runscript api.scripts.handle_security_scans"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0


